# backend/5g-agent-fw/config/5g-agent.service
[Unit]
Description=Secure Service Agent
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/5g-agent-fw
EnvironmentFile=/etc/environment
ExecStart=/opt/5g-agent-fw/startup.sh
Restart=always
RestartSec=15s

[Install]
WantedBy=multi-user.target
